# gc-image/crontab
20 * * * * for job in $(kubectl get pod | grep "scipion-tool-job" | grep "Completed" | cut -d" " -f1); do kubectl delete job "$job"; done
